# systemd/user/console.target
[Unit]
Description=Console target
Requires=dbus.socket
