# systemd/user/console.target
[Unit]
Description=Console interface target
Requires=default.target
After=default.target
; Requires=dbus.socket

; [Install]
; Alias=default.target
